# systemd/cloak-backup.service
[Unit]
Description=Cloak data backup (one-shot)
After=docker.service
Requires=docker.service
[Service]
Type=oneshot
ExecStart=/usr/local/bin/cloak-backup.sh
User=root
